# pmand/pmand.conf
[pmand]
pidfile=pmand.pid
logfile=pmand.log
directory=./

[program:sub.sh]
command=./sub.sh
logfile=sub.log
